# crontab
crontab -e
56 1 * * * /home/juanfelipeoz/Documentos/GPU-HunterBot/start.sh
